# dist/linux/scs.service
[Unit]
Description=SGX Caching Service

[Service]
Type=simple
User=scs
Group=scs
ExecStart=/usr/bin/scs run
ExecReload=/bin/kill -s HUP $MAINPID
TimeoutStartSec=0
Restart=on-failure
PermissionsStartOnly=true
RuntimeDirectory=scs
RuntimeDirectoryMode=0775

[Install]
WantedBy=multi-user.target
